# monitoring/systemd/zeek-relay.service
[Unit]
Description=Zeek conn.log relay to zeek_feed.log
After=network.target

[Service]
Type=simple
ExecStart=/home/ceo/honeypot/scripts/zeek_relay.sh
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
